# deployment/ansible/roles/configure-streamer/templates/telemetry-streamer.service.j2
[Unit]
Description=Telemetry Streaming Service
Documentation=https://github.com/your-org/telemetry-streaming
After=docker.service
Wants=docker.service

[Service]
Type=simple
User=telemetry
Group=telemetry
WorkingDirectory=/var/lib/telemetry

# Load environment variables
EnvironmentFile=/var/lib/telemetry/config/telemetry.env

# Service execution
ExecStart=/var/lib/telemetry/start-streamer.sh

# Service restart policy
Restart=always
RestartSec=10s
StartLimitInterval=600
StartLimitBurst=5

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security
PrivateTmp=yes
NoNewPrivileges=true

# Standard streams
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telemetry-streamer

[Install]
WantedBy=multi-user.target
